# infra/installer/templates/systemd/netx-core-service.service
[Unit]
Description=NetX Core Service (NestJS)
Documentation=https://github.com/CharlesBarreto1/netx
After=network-online.target postgresql.service redis-server.service rabbitmq-server.service
Wants=network-online.target postgresql.service redis-server.service rabbitmq-server.service
# Restart rate-limit fica em [Unit] em systemd 230+. Em [Service] dispara
# "Unknown key" warning no journal.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=netx
Group=netx
WorkingDirectory=/opt/netx/apps/core-service
EnvironmentFile=/etc/netx/.env
# Path profundo: como o tsconfig.json não fixa `rootDir` (removido pra
# permitir typecheck dos imports @netx/shared e @netx/config), o nest build
# calcula o common parent e o output reflete a estrutura completa do monorepo
# dentro de dist/. Por isso main.js fica em dist/apps/core-service/src/.
ExecStart=/usr/bin/node dist/apps/core-service/src/main.js

# Restart policy
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netx-core

# Security hardening
#
# NoNewPrivileges INTENCIONALMENTE OMITIDO. Quando true, o kernel ignora o
# bit setuid de /usr/bin/sudo dentro deste processo, o que QUEBRA o hook
# do NetworkEquipmentService que dispara sync-firewall.sh + sync-ntp.sh via
# `execFile('sudo', ['-n', script])`. Resultado seria: cadastrar NAS pela
# UI não abriria UFW automaticamente nem recarregaria FreeRADIUS — operador
# teria que rodar sync manualmente em cada cadastro.
#
# Trade-off: o user `netx` é unprivileged e o sudoers (/etc/sudoers.d/
# netx-infra-sync) restringe NOPASSWD APENAS a 4 comandos específicos
# (2 scripts + reload/restart freeradius). Sem ALL=ALL=NOPASSWD aberto,
# o blast radius de um RCE no backend permanece pequeno.
#
# TODO v1.1: migrar pra systemd Path activation (netx escreve flag file
# em /var/lib/netx/sync-requested, um .path unit detecta inotify e dispara
# um .service rodando como root). Aí pode voltar NoNewPrivileges=true.
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# IMPORTANTE: ProtectSystem=strict deixa todo / read-only no mount namespace
# DESTE cgroup — inclui processos escalados via sudo. Sem expor explicitamente
# /etc/chrony e /etc/ufw em ReadWritePaths, os scripts sync-ntp.sh e
# sync-firewall.sh falham silenciosamente ao tentar escrever:
#   sync-ntp:      /etc/chrony/conf.d/netx-allows.conf
#   sync-firewall: ufw escreve em /etc/ufw/user.rules
#
# /var/backups/netx: o backup MANUAL via UI (BackupsService) roda DENTRO deste
# processo e escreve o .dump direto via pg_dump. Sem expor o dir aqui, o
# pg_dump aborta com EROFS ("read-only file system") mesmo o dir existindo —
# o backup automático (netx-backup.timer) não sofre porque é um .service à
# parte com namespace próprio. Criado pelo installer (backups_install_dirs),
# então sempre existe; não precisa do prefixo `-`.
#
# CUIDADO: systemd RECUSA criar o namespace se algum path em ReadWritePaths
# NÃO EXISTIR — serviço fica em crash loop com status=226/NAMESPACE e
# "Failed to set up mount namespacing: /path: No such file or directory".
# UFW no Debian 13 NÃO USA /var/lib/ufw (estado fica em /etc/ufw/), então
# NÃO inclua aqui mesmo que outras distros usem.
#
# Prefixo `-` em ReadWritePaths torna o path opcional: systemd ignora se
# não existir em vez de derrubar o serviço. Defesa em profundidade pra
# diferenças entre distros (ex.: Ubuntu sem chrony default).
#
# Não adicione /etc/sudoers.d aqui — ele NÃO precisa ser writable em runtime
# (install.sh popula durante setup) e expor sudoers ao backend = vetor de RCE.
ReadWritePaths=/var/lib/netx /var/log/netx /var/backups/netx /opt/netx/apps/core-service/dist -/etc/chrony -/etc/ufw
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

# Limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
