# infra/installer/templates/systemd/netx-api-gateway.service
[Unit]
Description=NetX API Gateway (NestJS)
Documentation=https://github.com/CharlesBarreto1/netx
After=network-online.target netx-core-service.service
Wants=network-online.target
Requires=netx-core-service.service
# Rate-limit fica em [Unit] (systemd 230+). Em [Service] dispara warning.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=netx
Group=netx
WorkingDirectory=/opt/netx/apps/api-gateway
EnvironmentFile=/etc/netx/.env
# Path profundo (idem netx-core-service.service): tsconfig sem rootDir gera
# dist/apps/<app>/src/main.js em vez de dist/main.js.
ExecStart=/usr/bin/node dist/apps/api-gateway/src/main.js

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=netx-gateway

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/netx /var/log/netx /opt/netx/apps/api-gateway/dist
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
