# infra/installer/templates/systemd/netx-core-service.service
[Unit]
Description=NetX Core Service (NestJS)
Documentation=https://github.com/CharlesBarreto1/netx
After=network-online.target postgresql.service redis-server.service rabbitmq-server.service
Wants=network-online.target postgresql.service redis-server.service rabbitmq-server.service
# Restart rate-limit fica em [Unit] em systemd 230+. Em [Service] dispara
# "Unknown key" warning no journal.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=netx
Group=netx
WorkingDirectory=/opt/netx/apps/core-service
EnvironmentFile=/etc/netx/.env
# Path profundo: como o tsconfig.json não fixa `rootDir` (removido pra
# permitir typecheck dos imports @netx/shared e @netx/config), o nest build
# calcula o common parent e o output reflete a estrutura completa do monorepo
# dentro de dist/. Por isso main.js fica em dist/apps/core-service/src/.
ExecStart=/usr/bin/node dist/apps/core-service/src/main.js

# Restart policy
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netx-core

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/netx /var/log/netx /opt/netx/apps/core-service/dist
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

# Limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
